# infra/nitro/systemd/subly402-parent.service
[Unit]
Description=Subly402 Nitro parent relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=subly402
Group=subly402
WorkingDirectory=/var/lib/subly402
ExecStart=/opt/subly402/bin/start-parent.sh /etc/subly402/parent.env
Restart=always
RestartSec=5
LimitNOFILE=65536
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/subly402 /var/log/subly402 /etc/subly402 /opt/subly402

[Install]
WantedBy=multi-user.target
